# deploy/systemd/journal-sltp-poller.service.template
[Unit]
Description=Lingonberry Journal SL TP Poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_USER@
WorkingDirectory=@APP_DIR@
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-@APP_DIR@/.env
ExecStart=@APP_DIR@/.venv/bin/python3 jobs/sltp_poller.py continuous
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
