# deploy/systemd/journal-webapp.service.template
[Unit]
Description=Lingonberry Journal Web App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_USER@
WorkingDirectory=@APP_DIR@
Environment=PYTHONUNBUFFERED=1
Environment=WEBAPP_PORT=5000
EnvironmentFile=-@APP_DIR@/.env
ExecStart=@APP_DIR@/.venv/bin/gunicorn --workers 2 --threads 4 --timeout 120 --bind 0.0.0.0:${WEBAPP_PORT} webapp.app:app
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
